# develop/grunt.conf
[program:grunt]
directory=/vagrant/
command=/vagrant/node_modules/grunt-cli/bin/grunt watch
autostart=true
autorestart=true