# develop/jekyll.conf
[program:jekyll]
directory=/vagrant/
command=/usr/local/rvm/wrappers/adicu/bundle exec jekyll serve --watch
autostart=true
autorestart=true